# dot_config/systemd/user/auto_dark_mode.service
[Unit]
Description="Auto switch light/dark mode"

[Service]
ExecStart=%h/.config/waybar/scripts/dark-mode-manager.sh auto 8:00 18:00
